# transbordo.service
[Unit]
Description=Trasbordo CNN script modbus
After=network.target

[Service]
ExecStart=/home/jetson/Documents/APP/Transbordo/env/bin/python3.11  /home/jetson/Documents/APP/Transbordo/main.py
WorkingDirectory=/home/jetson/Documents/APP/Transbordo/
User=jetson
Restart=always

[Install]
WantedBy=multi-user.target
